# ops/systemd/ddclass-backend.service
[Unit]
Description=DD-CLASS FastAPI backend
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=ddclass
Group=ddclass
WorkingDirectory=/opt/dd-class/source
Environment=PYTHONUNBUFFERED=1
Environment=PORT=8001
Environment=GUNICORN_WORKERS=3
Environment=LOG_LEVEL=info
Environment=UPLOADS_DIR=/opt/dd-class/shared/uploads
EnvironmentFile=/opt/dd-class/shared/.env.production
ExecStartPre=/opt/dd-class/venv/bin/python -m app.bootstrap
ExecStart=/bin/bash -lc '/opt/dd-class/venv/bin/gunicorn --workers "${GUNICORN_WORKERS:-3}" --worker-class uvicorn.workers.UvicornWorker --bind 127.0.0.1:"${PORT:-8001}" --log-level "${LOG_LEVEL:-info}" --access-logfile - --error-logfile - app.main:app'
Restart=always
RestartSec=5
TimeoutStartSec=90
KillSignal=SIGQUIT
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
